# backend/launchd/com.postpunk.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.postpunk.api</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>cd /Users/ash/Desktop/N8tiveFlow/backend &amp;&amp; /usr/local/opt/node@20/bin/node server.mjs</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/ash/Desktop/N8tiveFlow/backend</string>

    <key>StandardOutPath</key>
    <string>/Users/ash/Desktop/N8tiveFlow/backend/api.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ash/Desktop/N8tiveFlow/backend/api.err.log</string>
  </dict>
</plist>
