# backend/systemd/postpunk-api.service
[Unit]
Description=PostPunk Backend API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/postpunk/backend
EnvironmentFile=/opt/postpunk/.env
ExecStart=/usr/bin/env node server.mjs
Restart=always
RestartSec=5
User=postpunk
Group=postpunk
StandardOutput=append:/opt/postpunk/backend/api.log
StandardError=append:/opt/postpunk/backend/api.err.log

[Install]
WantedBy=multi-user.target

